# supervisor/watch-and-reload-apache2.conf
[program:watch-and-reload-apache2]
command=/bin/bash -c "watch-apache2-conf.sh $DAV_SVN_CONF_PATH"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
